# novnc-itksnap/conf.d/itksnap.conf
[program:itksnap]
command=itksnap -g "%(ENV_IMAGE_PATH)s"
autorestart=true
