# itksnap/conf.d/itksnap.conf
[program:itksnap]
command=itksnap
autorestart=false
