# itksnap/conf.d/ratpoison.conf
[program:ratpoison]
command=x-window-manager
autorestart=false
